# partitions/__aosp_bootimgs/boot/initrd/etc/inittab
::sysinit:/etc/rc.local
::respawn:/bin/cttyhack /bin/login shell
